# crontab
# this is cron tab for user
*   *   *   *   *   sh /home/rahul/cronfile/exampl

