# crontab
# this is cron tab for user
20  17   *   *   *    /home/rahul/cronfile/extractor.sh

